# cloudbox/roundcube/files/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info
childlogdir=/var/log/supervisor
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisord.sock

[program:nginx]
user=roundcube
autostart=true
autorestart=true
command=dockerize -template /opt/roundcube/config/config.inc.php.tmpl:/opt/roundcube/config/config.inc.php nginx

[program:php7-fpm]
user=roundcube
autostart=true
autorestart=true
command=php-fpm7 -F
