# utilities/traefik/files/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info
childlogdir=/var/log/supervisor
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisord.sock

[program:traefik]
user=traefik
autostart=true
autorestart=true
stopasgroup=true
command=dockerize -template etc/traefik/traefik.conf.tmpl:/etc/traefik/traefik.conf /usr/local/bin/traefik -C /etc/traefik/traefik.conf
